# supervisord.conf
[supervisord]
nodaemon=true
user=user
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:web]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 7860
directory=/home/user/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/web.log
stderr_logfile=/var/log/supervisor/web.log
environment=HOME="/home/user"

[program:hermes-gateway]
command=bash /home/user/app/scripts/hermes_gateway.sh
directory=/home/user/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/gateway.log
stderr_logfile=/var/log/supervisor/gateway.log
environment=HOME="/home/user"

[program:hermes-dashboard]
command=hermes dashboard --host 127.0.0.1 --insecure
directory=/home/user
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/dashboard.log
stderr_logfile=/var/log/supervisor/dashboard.log
environment=HOME="/home/user",PATH="/home/user/.hermes/bin:/home/user/.local/bin:/usr/local/bin:/usr/bin:/bin",DASHBOARD_PORT="9119"
